# connection-check.service
[Unit]
# Execute connection-check 

[Service]
ExecStart=/PATH/TO/SCRIPT/connection-check.sh
Description=connection-check
