# connection-check.timer
[Unit]
Description=connection-check

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=connection-check.service

[Install]
WantedBy=basic.target
